# orangepi/cold-fusion-sight.service
[Unit]
Description=Cold Fusion Sight (camera + control web UI)
# `network.target` (any iface) instead of `network-online.target` — the
# latter waits for NetworkManager-wait-online which can stall 10-60s on
# the Orange Pi if any configured iface fails to come up. Sight only
# needs a kernel-level network to bind 0.0.0.0:8080; rio reachability
# is checked at runtime by the NT4 client.
After=network.target
Wants=network.target
# Don't ever give up: the default StartLimit kicks the service into
# 'failed' after 5 restarts in 10 seconds, which is the wrong default for
# a robot field service. We want it to keep coming back, period.
StartLimitIntervalSec=0

[Service]
Type=simple
# Installer rewrites these via `sed` at install time so the service runs as
# whichever user holds the install dir (typically `orangepi`).
User=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__INSTALL_DIR__/sight.env
# `--timeout-graceful-shutdown 3` gives uvicorn 3 seconds to drain SSE +
# MJPEG clients on SIGTERM before forcibly closing them. Without it,
# uvicorn defaults to a longer drain window and `systemctl restart`
# noticeably hangs in the pit.
ExecStart=__INSTALL_DIR__/.venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8080 --timeout-graceful-shutdown 3
Restart=always
RestartSec=2
# Hard cap how long systemd waits for SIGTERM to take effect before
# escalating to SIGKILL. The default is 90s — way too long for a service
# the operator restarts repeatedly during a session. 5s is plenty: the
# graceful-shutdown timeout above is 3s, and clients should have closed
# by then.
TimeoutStopSec=5
# v4l2 needs read access to /dev/video0; that's `video` group on Debian.
SupplementaryGroups=video
# Sane resource limits — the relay shouldn't ever balloon, but if it does
# we'd rather have systemd kick it than starve the rest of the Pi.
LimitNOFILE=4096
TasksMax=128
# Keep stdout/stderr in journalctl with a short tag for easy filtering
# (`journalctl -u cold-fusion-sight -f` shows live logs).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cold-fusion-sight

[Install]
WantedBy=multi-user.target
