# orangepi/cold-fusion-sight.service
[Unit]
Description=Cold Fusion Sight (camera + control web UI)
After=network-online.target
Wants=network-online.target
# Don't ever give up: the default StartLimit kicks the service into
# 'failed' after 5 restarts in 10 seconds, which is the wrong default for
# a robot field service. We want it to keep coming back, period.
StartLimitIntervalSec=0

[Service]
Type=simple
# Installer rewrites these via `sed` at install time so the service runs as
# whichever user holds the install dir (typically `orangepi`).
User=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__INSTALL_DIR__/sight.env
ExecStart=__INSTALL_DIR__/.venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
# v4l2 needs read access to /dev/video0; that's `video` group on Debian.
SupplementaryGroups=video
# Sane resource limits — the relay shouldn't ever balloon, but if it does
# we'd rather have systemd kick it than starve the rest of the Pi.
LimitNOFILE=4096
TasksMax=128
# Keep stdout/stderr in journalctl with a short tag for easy filtering
# (`journalctl -u cold-fusion-sight -f` shows live logs).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cold-fusion-sight

[Install]
WantedBy=multi-user.target
