# acuity/firmware/files/etc/systemd/system/acuity-dashboard.service
[Unit]
Description=Acuity (camera + control web UI)
# `network.target` (any iface) instead of `network-online.target` —
# we want the dashboard to come up even in AP mode, where there's no
# upstream link. AP-mode-only the dashboard reads "rio disconnected"
# but the operator can still see the camera feed.
After=network.target acuity-firstboot.service
Wants=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=acuity
WorkingDirectory=/opt/acuity/sight
EnvironmentFile=-/opt/acuity/sight/sight.env
ExecStart=/opt/acuity/.venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8080 --timeout-graceful-shutdown 3
Restart=always
RestartSec=2
TimeoutStopSec=5
SupplementaryGroups=video
LimitNOFILE=4096
TasksMax=128
StandardOutput=journal
StandardError=journal
SyslogIdentifier=acuity-dashboard
StateDirectory=acuity
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
